# systemd/freshrss-yt-favicons.service
[Unit]
Description=Refresh YouTube channel avatar favicons for FreshRSS

[Service]
Type=oneshot
ExecStart=/usr/local/bin/freshrss-yt-favicons.sh
