# systemd/freshrss-fetch.service
[Unit]
Description=Batched FreshRSS feed actualization

[Service]
Type=oneshot
ExecStart=/usr/local/bin/freshrss-fetch.sh
