# systemd/freshrss-fetch.timer
[Unit]
Description=Batched FreshRSS feed actualization (every 10 min)

[Timer]
OnCalendar=*:0/10
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
